# modules/e10pro/hosting/server/etc/cron.d/e10-main-hosting-server.conf
18  20    * * *   root    cd /var/www/data-sources/50842906798390 && e10-app cli-action --action=e10pro.hosting.server/master-certs-scan
